# dot_config/systemd/user/1password.service
[Unit]
Description=1Password Desktop Application
After=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/1password --silent
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target